# guest/supervisord.conf
[unix_http_server]
file=/run/hermes-box-supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/workspace/hermes-home/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/run/hermes-box-supervisord.pid
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/hermes-box-supervisor.sock

[program:sshd]
command=/usr/sbin/sshd -D -e
priority=10
autostart=true
autorestart=unexpected
startsecs=1
stdout_logfile=/workspace/hermes-home/logs/sshd.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
redirect_stderr=true
stopasgroup=true
killasgroup=true

[program:hermes]
command=/usr/local/bin/hermes gateway run --replace
user=hermes
directory=/workspace/work
environment=HOME="/home/hermes",HERMES_HOME="/workspace/hermes-home",PATH="/usr/local/bin:/usr/bin:/bin"
priority=20
autostart=true
autorestart=unexpected
startsecs=3
startretries=20
stdout_logfile=/workspace/hermes-home/logs/hermes-gateway.log
stdout_logfile_maxbytes=25MB
stdout_logfile_backups=5
redirect_stderr=true
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=30
